# deploy/banendev.service
[Unit]
Description=banen.dev (Next.js)
After=network.target

[Service]
Type=simple
# Pas User en WorkingDirectory aan naar jouw situatie:
User=www-data
WorkingDirectory=/var/www/banendev
Environment=NODE_ENV=production
# Shared VPS (all this owner's): :3000 Express site, :3001 gtmbanen, :3002 maakjobs.
# banen.dev uses :3003. Confirm free on the box:  sudo ss -ltnp | grep ':3003'
Environment=PORT=3003
Environment=NEXT_PUBLIC_SITE_URL=https://banen.dev
# Secrets (RESEND_API_KEY, ALERTS_FROM_EMAIL, SUBMISSIONS_TO_EMAIL, ADMIN_TOKEN,
# NEXT_PUBLIC_GA_ID, NEXT_PUBLIC_CLARITY_ID) - optional file, see .env.example:
EnvironmentFile=-/var/www/banendev/.env
# Roept de Next.js-server direct aan (robuuster dan via npm in systemd).
# --hostname 127.0.0.1 = alleen bereikbaar via nginx (niet publiek zoals 3001/3002):
ExecStart=/usr/bin/node /var/www/banendev/node_modules/next/dist/bin/next start --hostname 127.0.0.1 --port 3003
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
