# deploy/gtmbanen.service
[Unit]
Description=GTM Banen (Next.js)
After=network.target

[Service]
Type=simple
# Pas User en WorkingDirectory aan naar jouw situatie:
User=www-data
WorkingDirectory=/var/www/gtmbanen
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXT_PUBLIC_SITE_URL=https://gtmbanen.nl
# Roept de Next.js-server direct aan (robuuster dan via npm in systemd):
ExecStart=/usr/bin/node /var/www/gtmbanen/node_modules/next/dist/bin/next start --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
